# rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Starts app
# ==============================================================================
cd /opt/app || exit

exec yarn start
